# storage.norc/etc/udev/rules.d/95-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="<PATH-TO>/scripts/monitorhotplug.sh"
